# overlay/base/usr/lib/systemd/system/os-checkupdates.service
[Unit]
Description=Check for updates

[Service]
ExecStart=/usr/bin/os checkupdates --download
SuccessExitStatus=0 2
